# data_processing/crontab
# Run the data processing script once a day at midnight
0 0 * * * root python /app/main.py >> /var/log/cron.log 2>&1