# data_processing/crontab
# Run the data processing script every 5 minutes
*/1 * * * * python3 /app/main.py >> /var/log/cron.log 2>&1
