# data_processing/crontab
# Run the data processing script once a day at 1 PM
0 13 * * * root python3 /main.py >> /var/log/cron.log 2>&1
